# docker/debian-s6-java25/rootfs/etc/s6-overlay/s6-rc.d/network-node-autostart/up
#!/command/execlineb -P
/etc/s6-overlay/scripts/network-node-autostart
