# docker/ubi8-s6-java21/rootfs/etc/services.d/network-node/finish
#!/command/with-contenv bash
set -euo pipefail

exit_code="${1:-0}"
if [ "${exit_code}" -eq 0 ]; then
  if command -v s6-svc >/dev/null 2>&1; then
    s6-svc -d /run/service/network-node 2>/dev/null || true
    s6-svc -d /run/services/network-node 2>/dev/null || true
  fi
fi

exit 0
