# docker/ubi8-s6-java21/rootfs/etc/services.d/network-node/run
#!/command/with-contenv bash
set -euo pipefail
exec s6-setuidgid hedera /opt/hgcapp/services-hedera/HapiApp2.0/entrypoint.sh
