# docker/ubi8-s6-java25/rootfs/etc/s6-overlay/s6-rc.d/network-node-autostart/up
#!/command/with-contenv bash
set -euo pipefail

FLAG="/opt/hgcapp/services-hedera/HapiApp2.0/state/network-node.enabled"
if [ -f "${FLAG}" ]; then
  exec /command/s6-rc -u change network-node
fi
