# docker/ubi8-s6-java25/rootfs/etc/s6-overlay/s6-rc.d/network-node-log/run
#!/command/execlineb -P
s6-setuidgid hedera s6-log -b -- t /var/log/network-node
